# .config/systemd/user/backup.service
[Unit]
Description=Run a borg backup

[Service]
Type=simple
ExecStartPre=notify-send "Backup" "Starting daily borg backup"
ExecStart=/home/markus/.local/bin/backup
ExecStartPost=notify-send "Backup" "Daily borg backup done!"
