# .config/systemd/user/backup.timer
[Unit]
Description=Timer for weekly backups using borg

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
