# .config/systemd/user/bspwm.service
[Unit]
Description=Bspwm
BindsTo=graphical-session.target

[Service]
ExecStart=/usr/bin/bspwm
ExecStop=/usr/bin/bspc quit

