# .config/systemd/user/bspwm-session.target
[Unit]
Description=Bspwm session running
BindsTo=graphical-session.target

