# systemd/.config/systemd/user/backup.timer
[Unit]
Description=Timer for daily backups using borg

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
